# dist/a2text-voice.service.example
[Unit]
Description=a2text voice daemon (push-to-talk dictation)
Documentation=https://gitlab.com/xmsger/a2text

# graphical-session.target ensures we boot once a desktop session is active.
# pipewire/pipewire-pulse cover the audio capture backend selection inside
# the daemon — without these the depcheck at startup would land in a
# restart loop until PipeWire is up.
After=graphical-session.target pipewire.service pipewire-pulse.service
Wants=pipewire.service pipewire-pulse.service
PartOf=graphical-session.target

[Service]
# Type=notify is critical: systemd waits for sd_notify(READY=1) before
# marking the unit active, and the daemon's per-state STATUS= updates
# show up in `systemctl --user status a2text-voice`.
Type=notify
NotifyAccess=main

# IMPORTANT: --daemon forces the daemon-only entry point. Without it
# `a2text` (no flags) self-bootstraps, which under systemd would mean
# "if a daemon is already running, send Toggle and exit" — making this
# unit go inactive immediately. Always keep --daemon here.
ExecStart=/usr/local/bin/a2text --daemon

# systemd --user often has a stripped PATH. Spell it out so the daemon's
# subprocess lookups for pw-record / parecord / wl-copy succeed.
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Default config location for a per-user install. Adjust if you keep your
# config elsewhere; leave commented to use the default search path
# (./a2text.yaml or ./app/a2text.yaml relative to WorkingDirectory, which
# is unspecified under systemd → effectively unreliable).
Environment="A2TEXT_CONFIG=%h/.config/a2text/a2text.yaml"

# Restart on crash but not on clean exit (e.g. depcheck failure with
# missing required deps would loop forever otherwise).
Restart=on-failure
RestartSec=5

# Give the daemon time to drain in-flight transcription before SIGKILL.
TimeoutStopSec=20

# Disable core dumps. The daemon holds open /dev/input/event* fds and
# its heap transiently contains recently-read input_event packets
# (keycode + value pairs from foreign devices: TTY passwords, GPG
# passphrase prompts, USB-HID OTP tokens). A core dump captured after
# a panic would expose those bytes to anyone who can read
# /var/lib/systemd/coredump/. Forbid the dump outright.
LimitCORE=0

# --- Sandbox hardening ---
#
# The daemon needs read access to /dev/input/event* (hotkey backend),
# write access to /dev/uinput (autopaste), and the user's $XDG_*_HOME
# directories. Everything else is restricted. Each directive narrows the
# attack surface available to an attacker who gains RCE inside the
# daemon process — useful given the daemon already holds two kernel
# input handles (read keystrokes, write keystrokes) for its lifetime.
#
# `systemd-analyze security a2text-voice.service` will rate the unit
# after enabling these. Tune ReadWritePaths if you keep config/models
# outside the XDG defaults.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/a2text %h/.local/share/a2text %h/.cache/a2text
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@mount @swap @reboot @raw-io @cpu-emulation @obsolete @debug @clock @module
CapabilityBoundingSet=
AmbientCapabilities=
# Network access required for cloud STT providers (OpenAI, Deepgram) and
# go-whisper HTTP backend. Restrict to IPv4/IPv6 only — no AF_PACKET, no
# raw sockets. If you exclusively use whisper.cpp (local), tighten this
# to `RestrictAddressFamilies=AF_UNIX` to lock the daemon offline.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# /dev/uinput + /dev/input/event* must remain accessible. DeviceAllow is
# additive over the default empty set when DevicePolicy=closed.
DevicePolicy=closed
DeviceAllow=/dev/uinput rw
DeviceAllow=char-input r

[Install]
WantedBy=graphical-session.target
